# errbit/templates/default/erbbit.service.erb
[Unit]
Description=erbbit
Requires=network.target

[Service]
Type=simple
User=<%= node[:erbbit][:erbbituser] %>
Group=<%= node[:erbbit][:erbbitgroup] %>
WorkingDirectory= <%= node[:erbbit][:erbbitinstalldir] %>/httpdocs
ExecStart=/usr/bin/bash -lc 'bundle exec rails server -e production --bind <%= node[:erbbit][:host] %> --port <%= node[:erbbit][:port] %>'
TimeoutSec=30
RestartSec=15s
Restart=always

[Install]
WantedBy=multi-user.target

